# templates/default/etc/systemd/system/cq.service.erb
[Unit]
Description=<%= @daemon_name %>
After=network.target

[Service]
Type=simple
RemainAfterExit=yes
PIDFile=<%= @cq_home %>/crx-quickstart/conf/cq.pid
EnvironmentFile=<%= @conf_file %>
ExecStart=<%= @cq_home %>/crx-quickstart/bin/start \
  > <%= @cq_home %>/crx-quickstart/logs/start.log 2>&1
ExecStop=<%= @cq_home %>/crx-quickstart/bin/stop \
  > <%= @cq_home %>/crx-quickstart/logs/stop.log 2>&1
User=<%= @user %>
LimitNOFILE=<%= @fd_limit %>

[Install]
WantedBy=multi-user.target
